# com.scholar-dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.scholar-dashboard</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>-m</string>
        <string>http.server</string>
        <string>80</string>
        <string>--bind</string>
        <string>127.0.0.1</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/htnamus/All_Stuff/Programming_Stuff/ScholarDashboard</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/scholar-dashboard.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/scholar-dashboard.err</string>
</dict>
</plist>
